# C34351/WEPDemo/iperfclient.service
[Unit]
Description=IperfServer
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/local/bin/iperfclient.sh

[Install]
WantedBy=multi-user.target

